# claude-code-terminal/rootfs/etc/s6-overlay/s6-rc.d/ttyd/finish
#!/command/with-contenv bashio
# ==============================================================================
# Claude Code Terminal - ttyd cleanup on stop/crash
# ==============================================================================

declare exit_code=${1}

if [ "${exit_code}" -ne 0 ] && [ "${exit_code}" -ne 256 ]; then
    bashio::log.warning "ttyd exited with code ${exit_code}, restarting in 5 seconds..."
    sleep 5
fi
