# claude-code-terminal/rootfs/etc/profile.d/claude-env.sh
# Placeholder - overwritten by init-claude at startup
export XDG_CONFIG_HOME="/data/.config"
export XDG_CACHE_HOME="/data/.cache"
export HOME="/data"
export TMUX_TMPDIR="/tmp/tmux"
export TERM="xterm-256color"
