# claude-code-terminal/rootfs/etc/s6-overlay/s6-rc.d/init-claude/run
#!/command/with-contenv bashio
# ==============================================================================
# Claude Code Terminal - Initialization
# Sets up environment, auth, and tmux config from HA add-on options
# ==============================================================================

bashio::log.info "Initializing Claude Code Terminal..."

# ---------------------------------------------------------------------------
# Persistent directories
# ---------------------------------------------------------------------------
mkdir -p /data/.config/anthropic \
         /data/.config/claude \
         /data/.cache \
         /data/projects

# ---------------------------------------------------------------------------
# Read add-on options
# ---------------------------------------------------------------------------
ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
MODEL=$(bashio::config 'model')
YOLO_MODE=$(bashio::config 'yolo_mode')
TMUX_HISTORY=$(bashio::config 'tmux_history_limit')
SESSION_NAME=$(bashio::config 'session_name')

# ---------------------------------------------------------------------------
# Write API key for Claude Code
# ---------------------------------------------------------------------------
if bashio::var.has_value "${ANTHROPIC_API_KEY}"; then
    bashio::log.info "Configuring Anthropic API key..."
    echo "${ANTHROPIC_API_KEY}" > /data/.config/anthropic/api_key
    chmod 600 /data/.config/anthropic/api_key
else
    bashio::log.warning "No Anthropic API key configured! Claude Code will prompt for auth."
fi

# ---------------------------------------------------------------------------
# Configure Claude Code settings
# ---------------------------------------------------------------------------
if bashio::var.true "${YOLO_MODE}"; then
    bashio::log.info "YOLO mode enabled - Claude will auto-approve tool use"
    cat > /data/.config/claude/settings.json <<EOF
{
  "permissions": {
    "allow": ["Bash(*)", "Read(*)", "Write(*)", "Edit(*)", "Glob(*)", "Grep(*)"],
    "deny": []
  }
}
EOF
fi

# ---------------------------------------------------------------------------
# Write tmux configuration
# ---------------------------------------------------------------------------
bashio::log.info "Configuring tmux (scrollback: ${TMUX_HISTORY})..."
cat > /data/.tmux.conf <<EOF
# Claude Code Terminal - tmux configuration
set-option -g history-limit ${TMUX_HISTORY}
set-option -g mouse on
set-option -g default-terminal "screen-256color"
set-option -g status on
set-option -g status-position bottom
set-option -g status-style "bg=#1a1b26,fg=#a9b1d6"
set-option -g status-left " #[fg=#7aa2f7,bold]Claude Code #[fg=#565f89]| "
set-option -g status-right "#[fg=#565f89]Session: #[fg=#9ece6a]#{session_name} "
set-option -g status-left-length 30
set-option -g status-right-length 40

# UTF-8 support
set-option -g -q utf8 on

# Don't automatically rename windows
set-option -g allow-rename off

# Start window numbering at 1
set-option -g base-index 1

# Reduce escape time for faster key response
set-option -sg escape-time 10

# Increase display time for messages
set-option -g display-time 3000
EOF

# ---------------------------------------------------------------------------
# Write environment file sourced by all scripts
# ---------------------------------------------------------------------------
cat > /etc/profile.d/claude-env.sh <<EOF
export XDG_CONFIG_HOME="/data/.config"
export XDG_CACHE_HOME="/data/.cache"
export HOME="/data"
export ANTHROPIC_API_KEY="${ANTHROPIC_API_KEY}"
export CLAUDE_MODEL="${MODEL}"
export SESSION_NAME="${SESSION_NAME}"
export TMUX_TMPDIR="/tmp/tmux"
export TERM="xterm-256color"

# Ensure tmux socket dir exists
mkdir -p /tmp/tmux
chmod 700 /tmp/tmux
EOF

chmod 644 /etc/profile.d/claude-env.sh

bashio::log.info "Claude Code Terminal initialization complete."
