# claude-code-terminal/rootfs/etc/services.d/ttyd/finish
#!/usr/bin/env bash
# ==============================================================================
# Claude Code Terminal - ttyd cleanup on stop/crash
# ==============================================================================

declare exit_code=${1}

if [ "${exit_code}" -ne 0 ] && [ "${exit_code}" -ne 256 ]; then
    sleep 5
fi
