# claude-code-terminal/rootfs/etc/services.d/ttyd/run
#!/command/with-contenv bashio
# ==============================================================================
# Claude Code Terminal - ttyd web terminal server
# ==============================================================================

source /etc/profile.d/claude-env.sh

declare ingress_port
ingress_port=$(bashio::addon.ingress_port)

bashio::log.info "Starting ttyd on port ${ingress_port}..."

exec ttyd \
    --interface 0.0.0.0 \
    --port "${ingress_port}" \
    --writable \
    --reconnect 10 \
    --ping-interval 30 \
    --base-path "$(bashio::addon.ingress_entry)" \
    /usr/bin/tmux-session.sh
