# claude-code-terminal/rootfs/etc/services.d/ttyd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "=== Claude Code Terminal starting ==="

# -------------------------------------------------------
# Persistent storage setup
# -------------------------------------------------------
PERSIST_DIR=/data/claude-persist
mkdir -p "${PERSIST_DIR}/config" "${PERSIST_DIR}/cache" /data/projects

# -------------------------------------------------------
# Read add-on options via bashio
# -------------------------------------------------------
API_KEY=$(bashio::config 'anthropic_api_key' '')
MODEL=$(bashio::config 'model' 'sonnet')
YOLO_MODE=$(bashio::config 'yolo_mode' 'false')
TMUX_HISTORY=$(bashio::config 'tmux_history_limit' '50000')
SESSION_NAME=$(bashio::config 'session_name' 'claude')

# -------------------------------------------------------
# API key
# -------------------------------------------------------
if bashio::var.has_value "${API_KEY}"; then
    export ANTHROPIC_API_KEY="${API_KEY}"
    bashio::log.info "Anthropic API key configured"
else
    bashio::log.warning "No Anthropic API key set - Claude will prompt for auth"
fi

# -------------------------------------------------------
# Symlink Claude config to persistent storage
# -------------------------------------------------------
if [ ! -L /root/.claude ]; then
    rm -rf /root/.claude
    ln -s "${PERSIST_DIR}" /root/.claude
fi
mkdir -p /root/.config
if [ ! -L /root/.config/claude-code ]; then
    rm -rf /root/.config/claude-code
    ln -s "${PERSIST_DIR}/config" /root/.config/claude-code
fi

# -------------------------------------------------------
# Update tmux scrollback from options
# -------------------------------------------------------
sed -i "s/history-limit.*/history-limit ${TMUX_HISTORY}/" /root/.tmux.conf

# -------------------------------------------------------
# YOLO mode
# -------------------------------------------------------
if bashio::var.true "${YOLO_MODE}"; then
    bashio::log.info "YOLO mode enabled"
    echo '{"permissions":{"allow":["Bash(*)","Read(*)","Write(*)","Edit(*)","Glob(*)","Grep(*)"],"deny":[]}}' \
        > "${PERSIST_DIR}/settings.json"
fi

# -------------------------------------------------------
# Write environment file for tmux sessions to source
# -------------------------------------------------------
ENV_FILE=/data/claude-env.sh
cat > "${ENV_FILE}" <<ENVEOF
export TERM=xterm-256color
export LANG=C.UTF-8
export SHELL=/bin/bash
export CLAUDE_MODEL="${MODEL}"
export SESSION_NAME="${SESSION_NAME}"
export PATH="/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/root/.npm-global/bin:\${PATH}"
ENVEOF

if bashio::var.has_value "${API_KEY}"; then
    echo "export ANTHROPIC_API_KEY=\"${API_KEY}\"" >> "${ENV_FILE}"
fi

if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\"" >> "${ENV_FILE}"
fi

# -------------------------------------------------------
# Set default shell for tmux
# -------------------------------------------------------
export SHELL=/bin/bash
echo "set-option -g default-shell /bin/bash" >> /root/.tmux.conf
echo "set-option -g default-command /bin/bash" >> /root/.tmux.conf

# -------------------------------------------------------
# Home Assistant MCP server
# -------------------------------------------------------
if bashio::var.has_value "${SUPERVISOR_TOKEN:-}"; then
    bashio::log.info "Configuring Home Assistant MCP server..."
    claude mcp remove homeassistant -s user 2>/dev/null || true
    claude mcp add-json homeassistant \
        "{\"type\":\"http\",\"url\":\"http://supervisor/core/api/mcp\",\"headers\":{\"Authorization\":\"Bearer ${SUPERVISOR_TOKEN}\"}}" \
        -s user 2>/dev/null \
        && bashio::log.info "HA MCP configured" \
        || bashio::log.warning "HA MCP configuration failed (non-fatal)"
fi

# -------------------------------------------------------
# Create a bashrc that sources the env file
# -------------------------------------------------------
cat > /root/.bashrc <<'RCEOF'
[ -f /data/claude-env.sh ] && source /data/claude-env.sh
alias c="claude"
alias cc="claude --continue"
RCEOF

# -------------------------------------------------------
# Create startup script that launches Claude in tmux
# -------------------------------------------------------
CLAUDE_ARGS=""
if [ "${MODEL}" != "default" ] && [ -n "${MODEL}" ]; then
    CLAUDE_ARGS="--model ${MODEL}"
fi

cat > /data/claude-start.sh <<STARTEOF
#!/bin/bash
source /data/claude-env.sh
cd /data/projects 2>/dev/null || cd /root
echo ""
echo "  Claude Code Terminal for Home Assistant"
echo "  Session: ${SESSION_NAME} | Model: ${MODEL}"
echo "  Close this tab anytime - your session persists!"
echo ""
exec claude ${CLAUDE_ARGS}
STARTEOF
chmod +x /data/claude-start.sh

# -------------------------------------------------------
# Get the dynamic ingress port from the Supervisor
# -------------------------------------------------------
ingress_port=$(bashio::addon.ingress_port)
bashio::log.info "Starting ttyd on ingress port ${ingress_port}"

# -------------------------------------------------------
# Launch ttyd bound to the ingress port, wrapping tmux
# -------------------------------------------------------
cd /data/projects || cd /root

exec ttyd \
    -i 0.0.0.0 \
    -p "${ingress_port}" \
    -W \
    -t fontSize=14 \
    -t fontFamily="Monaco,Consolas,monospace" \
    /bin/bash -c "exec tmux new-session -A -s ${SESSION_NAME} /data/claude-start.sh"
